# deploy/nexora-monitor.service
[Unit]
Description=Nexora production health monitor
Wants=docker.service
After=docker.service network-online.target

[Service]
Type=oneshot
User=nexora
Group=nexora
WorkingDirectory=/opt/nexora/deploy
Environment=ENV_FILE=/opt/nexora/deploy/runtime/.env.production
Environment=COMPOSE_FILE=/opt/nexora/deploy/docker-compose.prod.yml
Environment=STATE_FILE=/opt/nexora/deploy/runtime/monitor-state
ExecStart=/opt/nexora/deploy/monitor-stack.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/nexora/deploy/runtime

